# launchd/com.rayjohnson.mcp-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rayjohnson.mcp-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ray.johnson/src/rayjohnson/mcp-gateway/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ray.johnson/src/rayjohnson/mcp-gateway</string>

    <!-- Restart on failure but not on intentional stop -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Prevent rapid restart loops on persistent failure -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ray.johnson/Library/Logs/mcp-gateway/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ray.johnson/Library/Logs/mcp-gateway/stderr.log</string>

    <!-- launchd runs with a minimal environment — set what start.sh needs -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ray.johnson</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
